# binaries/toolchain/usr/mipsel-ingenic-linux-uclibc/sysroot/lib/udev/rules.d/70-infrared.rules
# Automatically load the proper keymaps after the Remote Controller device
# creation.  The keycode tables rules should be at /etc/rc_maps.cfg

ACTION=="add", SUBSYSTEM=="rc", RUN+="/usr/bin/ir-keytable -a /etc/rc_maps.cfg -s $name"
